# docker/grafana-arm/supervisord.conf
[supervisord]
nodaemon=true

[program:grafana]
user=grafana
command=/run.sh

[program:sshd]
command=/usr/sbin/sshd -D
